# src/stemcell_builder/stages/system_azure_wala/assets/runit/waagent/run
#!/bin/bash
set -e

exec nice -n -10 /usr/sbin/waagent -daemon
