# stemcell_builder/stages/mcf_dirs/assets/tty1.conf
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn

exec /sbin/getty -n -i -l /var/vcap/micro/bin/microconsole -8 38400 tty1 -8 38400 tty1
