# stemcell_builder/stages/bosh_audit/assets/auditd_upstart.conf
# auditd - system audit daemon

description	"system audit daemon"

start on filesystem
stop on shutdown

script
    /etc/init.d/auditd start
end script
