# agent/misc/stemcell/instance/runit/agent/run
#!/bin/bash
exec 2>&1
exec /var/vmc/bosh/agent/bin/agent -c
