# stemcell_builder/stages/base_centos_packages/assets/runit.service
[Unit]
Description=runit
After=network.target

[Service]
TimeoutStartSec=0
ExecStart=/sbin/runsvdir-start

[Install]
WantedBy=multi-user.target
