# src/stemcell_builder/stages/base_warden/assets/runit/cron/run
#!/bin/bash

set -e

PATH=/usr/local/bin:/usr/local/sbin:/bin:/sbin:/usr/bin:/usr/sbin:/usr/X11R6/bin

exec cron -f
