# stemcell_builder/stages/bosh_go_agent/assets/runit/agent/run
#!/bin/bash
set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1

cd /var/vcap/bosh

exec nice -n -10 /var/vcap/bosh/bin/bosh-agent -P $(cat /var/vcap/bosh/etc/operating_system) -C /var/vcap/bosh/agent.json
