# stemcell_builder/stages/base_warden/assets/runit/rsyslog/run
#!/bin/bash
#
# Copyright (c) 2009-2012 VMware, Inc.

set -e

. /etc/default/rsyslog
exec /usr/sbin/rsyslogd -n $RSYSLOGD_OPTIONS
