# stemcell_builder/stages/rsyslog_config/assets/rsyslog_upstart.conf
# rsyslog - system logging daemon
#
# rsyslog is an enhanced multi-threaded replacement for the traditional
# syslog daemon, logging messages from applications

description	"system logging daemon"

start on filesystem
stop on runlevel [06]

expect fork
respawn

script
    . /etc/default/rsyslog
    exec /usr/sbin/rsyslogd $RSYSLOGD_OPTIONS
end script