# bosh_warden_cpi/init/bosh.conf
start on runlevel [2345]
stop on runlevel [!2345]

script
  start bosh-nats
  start bosh-redis
  start bosh-blobstore
  start bosh-warden
  start bosh-director
  start bosh-worker
end script
