# bosh_warden_cpi/init/bosh-redis.conf
start on runlevel [345]
stop on runlevel [!345]

script
  exec redis-server /usr/local/etc/redis.conf \
  1>>/var/log/bosh/redis.out.log \
  2>>/var/log/bosh/redis.err.log
end script
