# bosh_warden_cpi/init/bosh-nats.conf
start on runlevel [345]
stop on runlevel [!345]

script
  exec nats-server -a 0.0.0.0
end script
