# vmbuilder/templates/micro/instance/runit/agent/run
#!/bin/bash
exec 2>&1
exec /var/vcap/bosh/agent/bin/agent -c
