# bosh_agent/misc/stemcell/build2/stages/mcf_agent/assets/runit/agent/run
#!/bin/bash
#
# Copyright (c) 2009-2012 VMware, Inc.

set -e

export PATH=/var/vcap/bosh/bin:$PATH
cd /var/vcap/bosh/bosh/bosh_agent
exec 2>&1
exec bundle exec bin/bosh_agent -c -I $(cat /etc/infrastructure) -n http://vcap:vcap@localhost:6969 -P microcloud
