# src/stemcell_builder/stages/bosh_go_agent/assets/runit/monit/run
#!/bin/bash
set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1

exec nice -n -10 /var/vcap/bosh/bin/monit -I -c /var/vcap/bosh/etc/monitrc