# stemcell_builder/stages/system_vsphere_cdrom/assets/95-bosh-cdrom.rules
KERNEL=="sr0", ACTION=="change", RUN+="/etc/udev/rules.d/ready_cdrom.sh"