# stemcell_builder/stages/bosh_agent/assets/runit/agent/run
#!/bin/bash
#
# Copyright (c) 2009-2012 VMware, Inc.

set -e

export PATH=/var/vcap/bosh/bin:$PATH
exec 2>&1

exec nice -n -10 /var/vcap/bosh/bin/bosh_agent -c -I $(cat /var/vcap/bosh/etc/infrastructure) -P $(cat /var/vcap/bosh/etc/operating_system)
